# templates/default/nagios-nrpe-server.erb
#!/bin/sh
#
# Autogenerated by Chef.
#
# chkconfig: 2345 80 30
# description: nrpe is a daemon for a remote nagios server, \
#              running nagios plugins on this host.
# processname: <%= node['nagios']['nrpe']['service_name'] %>
# config: <%= node['nagios']['nrpe']['conf_dir'] %>/nrpe.cfg

<% if ['rhel', 'fedora'].include?(node['platform_family']) -%>
# Source function library
if [ -f /etc/rc.d/init.d/functions ]; then
. /etc/rc.d/init.d/functions
elif [ -f /etc/init.d/functions ]; then
. /etc/init.d/functions
elif [ -f /etc/rc.d/functions ]; then
. /etc/rc.d/functions
fi


# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

<% end -%>
NrpeBin=/usr/bin/nrpe
NrpeCfg=<%= node['nagios']['nrpe']['conf_dir'] %>/nrpe.cfg

case "$1" in
start)
  # Start daemon.
  echo -n "Starting <%= node['nagios']['nrpe']['service_name'] %>: "
  <% if ['rhel', 'fedora'].include?(node['platform_family']) -%>daemon $NrpeBin -c $NrpeCfg -d<% else -%>start-stop-daemon --start --quiet --exec $NrpeBin -- -c $NrpeCfg -d<% end -%>
  echo
  ;;
stop)
  # Stop daemon.
  echo -n "Shutting down <%= node['nagios']['nrpe']['service_name'] %>: "
  <% if ['rhel', 'fedora'].include?(node['platform_family']) -%>killproc <%= node['nagios']['nrpe']['service_name'] %><% else -%>start-stop-daemon --stop --quiet --exec $NrpeBin<% end -%>
  echo
  ;;
restart)
  $0 stop
  sleep 1
  $0 start
  ;;
*)
  echo "Usage: <%= node['nagios']['nrpe']['service_name'] %> {start|stop|restart|status}"
  exit 1
esac

exit 0
